# services.d/yggdrasil/run
#!/bin/sh
exec 2>&1
exec /usr/bin/yggdrasil -useconffile /etc/yggdrasil.conf
